# systemd_units/minecraft-backup@.timer
[Unit]
Description=Schedule Minecraft Server Backup for %i

[Timer]
OnCalendar=daily
Persistent=yes

[Install]
WantedBy=timers.target

